# ansible/centos-gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{repo_dir}}/gunicorn.sock

[Install]
WantedBy=sockets.target
